# src/modules/tilty/filesystem/root/etc/systemd/system/tilty-watcher.path
[Path]
Unit=tilty-watcher.service
PathChanged=/etc/tilty/tilty.ini

[Install]
WantedBy=multi-user.target
